# systemd/borg_blair/borg_blair.service
[Unit]
Description=borg backup

[Service]
Type=simple
ExecStart=/home/blair/Projects/borg/blair/backup.sh

[Install]
WantedBy=default.target
