# systemd/photoprism/pp_start.timer
[Unit]
Description=start photoprism on boot

[Timer]
OnBootSec=1 m 
Unit=pp_start.service

[Install]
WantedBy=timers.target
