# systemd/transmission/transmission.service
[Unit]
Description=start transmission daemon

[Service]
Type=simple
ExecStart=transmission-daemon

[Install]
WantedBy=default.target




