# systemd/transmission/transmission.timer
[Unit]
Description=start transmission

[Timer]
OnBootSec=5min 
Unit=transmission.service

[Install]
WantedBy=timers.target
