# systemd/dropsync/dropsync.service
[Unit]
Description= sync to dropbox

[Service]
Type=simple
ExecStart=/home/blair/Dropbox/bin/dropsync

[Install]
WantedBy=default.target
