# systemd/borg_media/borg_media.service
[Unit]
Description=borg media backup

[Service]
Type=simple
ExecStart=/home/blair/Projects/borg/media/backup.sh

[Install]
WantedBy=default.target
