# systemd/media_backup/media_backup.service
[Unit]
Description= Backup media library

[Service]
Type=simple
ExecStart=/home/blair/cronjobs/media_backup

[Install]
WantedBy=default.target
